# camera/g500-camera.service
[Unit]
Description=G500 FPV Camera Stream
After=network.target

[Service]
Type=simple
User=yze
WorkingDirectory=/home/yze/g500_rc_car
ExecStart=/usr/bin/python3 /home/yze/g500_rc_car/camera/webcam_stream.py
Restart=always
RestartSec=3

# Environment overrides — edit these if your setup differs
Environment=CAMERA_INDEX=1
Environment=STREAM_WIDTH=640
Environment=STREAM_HEIGHT=480
Environment=JPEG_QUALITY=70
Environment=TARGET_FPS=25
Environment=STREAM_PORT=8080
Environment=STREAM_HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
